# deploy/systemd/netsentinel-ultra-backend.service
[Unit]
Description=Netsentinel Ultra Core Backend
After=network.target

[Service]
User=%i
WorkingDirectory=/home/%i/Netsentinel-ultra-core/backend
ExecStart=/home/%i/Netsentinel-ultra-core/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 5000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
